# serial_port_server/cros_console.service
# Systemd service file for Thermo CRS C500 serial port server
# Install to ~/.config/systemd/user/
# Enable with systemctl --user enable cros_console
#
# Connects CROS console serial port to TCP port 11010.
# The listener allows multiple connections, any can write to terminal and
# all see the replies.

[Unit]
Description=CROS Console Serial Port

[Service]
ExecStart=/usr/bin/ncat -l -k localhost 11010
ExecStartPost=/bin/stty -F /dev/ttyUSB2 57600 cs8 -cstopb -parenb -echo -echoe -echok -icanon -isig -opost -onlcr onocr -icrnl -ixon -ixoff crtscts
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/ttyUSB2

[Install]
WantedBy=default.target
